# personal/settings/.config/autostart/network-manager-applet.desktop
[Desktop Entry]
Name=Network Manager Applet
GenericName=Network Manager Applet
Exec=nm-applet
Terminal=false
Type=Application
StartupNotify=false
